# helloJob.cron
* * * * * echo "TEST 2" >> /cron_demo/log_file.txt
